# orchestral/actors/automail2note_legacy/workspace_automation/mail_to_calendar.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lazyingart.mailcalendar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/lachlan/.openclaw/workspace/.venv/bin/python</string>
        <string>/Users/lachlan/.openclaw/workspace/automation/mail_to_calendar.py</string>
    </array>

    <key>StartInterval</key>
    <integer>43200</integer>

    <key>WorkingDirectory</key>
    <string>/Users/lachlan/.openclaw/workspace</string>

    <key>StandardOutPath</key>
    <string>/Users/lachlan/.openclaw/workspace/logs/mail_to_calendar.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lachlan/.openclaw/workspace/logs/mail_to_calendar.launchd.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
